# docker/config/supervisor/config.d/websockets.conf
[program:websockets]
directory=/var/www/
command=php artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=root
stdout_logfile=/var/www/storage/logs/websockets.log
